# sources/etc/systemd/system/docker.service.d/ngs.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon -H fd:// --insecure-registry docker.app.s
